# apps/rclone/root/etc/s6-overlay/s6-rc.d/rcloneui/run
#!/command/with-contenv bash
/rcloneui.sh