# apps/rclone/root/etc/s6-overlay/s6-rc.d/mounter/run
#!/command/with-contenv bash
/mounter.sh
