# apps/plex/root/etc/services.d/s6-cron/run
#!/command/with-contenv sh
exec s6-cronb -l 5 -c /etc/s6-cron/crontabs
