# apps/rclonebrowser/rootfs/etc/services.d/rclonefm/run
#!/usr/bin/with-contenv sh

echo "[$(basename "$(pwd)")] starting..."

/usr/bin/rclone rcd \
    --rc-web-gui \
    --rc-addr=0.0.0.0:5573 \
    --rc-web-gui-no-open-browser \
    --rc-no-auth \
    --transfers=1 \
    --bwlimit=10M \
    --multi-thread-streams=1 \
    /var/lib/rclonefm