# apps/rclone/rootfs/etc/s6-overlay/s6-rc.d/rclonewebdav/run
#!/command/with-contenv sh

s6-setuidgid ${PUID}:${PGID} \
/usr/bin/rclone serve webdav \
    --config=/config/rclone.conf \
    --addr :5574 \
    --user ${WEBDAV_USERNAME} \
    --pass "${WEBDAV_PASSWORD}" \
    --transfers=1 \
    --bwlimit=10M \
    --multi-thread-streams=1 \
    /storage