# apps/filebrowser/root/etc/s6-overlay/s6-rc.d/filebrowser/run
#!/command/with-contenv bash
cd /tmp
/filebrowser
