# apps/plex/root/etc/s6-overlay/s6-rc.d/mount-checker/run
#!/command/with-contenv bash
/mount-checker.sh