# apps/plex/root/etc/s6-overlay/s6-rc.d/transcode-killer/run
#!/command/with-contenv bash
python3 /transcode-killer.py
