# apps/filebrowser/root/etc/s6-overlay/s6-rc.d/elfterm/run
#!/command/with-contenv bash

# Make elfterm's config accessible in filebrowser
ln -s /home/.elfie/.config /storage/config/elfterm

# Launch elfterm
/elfterm.sh