# apps/filebrowser/root/etc/s6-overlay/s6-rc.d/elfterm/run
#!/command/with-contenv
/elfterm.sh bash