# apps/filebrowser/root/etc/s6-overlay/s6-rc.d/elfterm/run
#!/command/with-contenv bash

ln -s /storage/config/plextraktsync /home/elfie/.config

# Launch elfterm
/elfterm.sh