# apps/gotosocial/root/etc/s6-overlay/s6-rc.d/gotosocial/run
#!/command/with-contenv ash
cd /gotosocial
./gotosocial server start
