# apps/bluesky-pds/root/etc/s6-overlay/s6-rc.d/pds/run
#!/command/with-contenv ash
cd /app
node --enable-source-maps index.js
