# apps/cli_debrid/root/etc/s6-overlay/s6-rc.d/secondary/run
#!/command/with-contenv bash
cd /app/cli_battery
python /app/cli_battery/main.py