# apps/rclone/root/etc/s6-overlay/s6-rc.d/webdav-local/run
#!/command/with-contenv bash
/webdav-local.sh