# apps/cli_debrid/root/etc/s6-overlay/s6-rc.d/primary/run
#!/command/with-contenv bash
cd /app
python /app/main.py