# apps/plex/root/etc/s6-overlay/s6-rc.d/plex/run
#!/command/with-contenv bash
/entrypoint.sh