# apps/bluesky-pds/root/etc/s6-overlay/s6-rc.d/pds-cli/run
#!/command/with-contenv ash

/pds-cli.sh