# apps/bluesky-pds/root/etc/s6-overlay/s6-rc.d/pdsadmin/run
#!/command/with-contenv ash

/pdsadmin.sh